# etc/profile.d/artisan-aliases.sh
alias a="php artisan"
alias aj:s="php artisan jwt:secret"
alias ak:g="php artisan key:generate"
alias amg="php artisan migrate"
alias amg:f="php artisan migrate:fresh"
alias amg:i="php artisan migrate:install"
alias amg:rf="php artisan migrate:refresh"
alias amg:rb="php artisan migrate:rollback"
alias amg:rs="php artisan migrate:reset"
alias amg:s="php artisan migrate:status"
alias amk:a="php artisan make:auth"
alias amk:c="php artisan make:controller"
alias amk:cmd="php artisan make:command"
alias amk:e="php artisan make:event"
alias amk:f="php artisan make:factory"
alias amk:j="php artisan make:job"
alias amk:l="php artisan make:listener"
alias amk:m="php artisan make:model"
alias amk:mg="php artisan make:migration"
alias amk:ml="php artisan make:mail"
alias amk:mw="php artisan make:middleware"
alias amk:pl="php artisan make:policy"
alias amk:pr="php artisan make:provider"
alias amk:r="php artisan make:request"
alias amk:rs="php artisan make:resource"
alias amk:t="php artisan make:test"
alias ar:ch="php artisan route:cache"
alias ar:cl="php artisan route:clear"
alias ar:l="php artisan route:list"
alias av:p="php artisan vendor:publish"
alias avw:cl="php artisan view:clear"
