# etc/profile.d/artisan-aliases.sh
alias a="php artisan"

# auth
alias aa="php artisan auth"
alias aa:cr="php artisan auth:clear-resets"

# cache
alias ac="php artisan cache"
alias ac:c="php artisan cache:clear"
alias ac:f="php artisan cache:forget"
alias ac:t="php artisan cache:table"

# clear
alias acc="php artisan clear-compiled"

# config
alias acf="php artisan config"
alias acf:c="php artisan config:cache"
alias acf:cl="php artisan config:clear"

# db
alias adb="php artisan db"
alias adb:s="php artisan db:seed"
alias adb:w="php artisan db:wipe" 

# debugbar
alias adg="php artisan debugbar"
alias adg:c="php artisan debugbar:clear"

# down
alias ad="php artisan down"

# event
alias ae="php artisan event"
alias ae:c="php artisan event:cache"
alias ae:cl="php artisan event:clear"
alias ae:g="php artisan event:generate"
alias ae:l="php artisan event:list"

# env
alias aen="php artisan env"

# help
alias ah="php artisan help"

# inspire
alias ai="php artisan inspire"

# jwt
alias aj="php artisan jwt"
alias aj:s="php artisan jwt:secret"

# key
alias ak="php artisan key"
alias ak:g="php artisan key:generate"

# list
alias al="php artisan list"

# make
alias am="php artisan make"
alias am:a="php artisan make:auth"
alias am:c="php artisan make:controller"
alias am:cr="php artisan make:controller --resource"
alias am:ch="php artisan make:channel"
alias am:cm="php artisan make:command"
alias am:e="php artisan make:event"
alias am:ex="php artisan make:exception"
alias am:f="php artisan make:factory"
alias am:j="php artisan make:job"
alias am:l="php artisan make:listener"
alias am:m="php artisan make:model"
alias am:mg="php artisan make:migration"
alias am:ml="php artisan make:mail"
alias am:mw="php artisan make:middleware"
alias am:n="php artisan make:notification"
alias am:o="php artisan make:observer"
alias am:pl="php artisan make:policy"
alias am:pr="php artisan make:provider"
alias am:rq="php artisan make:request"
alias am:rs="php artisan make:resource"
alias am:rl="php artisan make:rule"
alias am:s="php artisan make:seeder"
alias am:t="php artisan make:test"

# migrate
alias amg="php artisan migrate"
alias amgs="php artisan migrate --seed"
alias amg:f="php artisan migrate:fresh"
alias amg:i="php artisan migrate:install"
alias amg:rf="php artisan migrate:refresh"
alias amg:rs="php artisan migrate:reset"
alias amg:rb="php artisan migrate:rollback"
alias amg:s="php artisan migrate:status"

# notification
alias an="php artisan notification:table"
alias an:t="php artisan notification:table"

# optimize
alias ao="php artisan optimize"
alias ao:c="php artisan optimize:clear"

# package
alias ap="php artisan package"
alias ap:d="php artisan package:discover"

# preset
alias apr="php artisan preset"

# queue
alias aq="php artisan queue"
alias aq:f="php artisan queue:failed"
alias aq:ft="php artisan queue:failed-table"
alias aq:fl="php artisan queue:flush"
alias aq:fg="php artisan queue:forget"
alias aq:l="php artisan queue:listen"
alias aq:rs="php artisan queue:restart"
alias aq:rt="php artisan queue:retry"
alias aq:t="php artisan queue:table"
alias aq:w="php artisan queue:work"

# route
alias ar="php artisan route"
alias ar:c="php artisan route:cache"
alias ar:cl="php artisan route:clear"
alias ar:l="php artisan route:list"

# schedule
alias asc="php artisan schedule"
alias asc:r="php artisan schedule:run"

# serve
alias as="php artisan serve"

# session
alias ase="php artisan session"
alias ase:t="php artisan session:table"

# storage
alias ast="php artisan storage"
alias ast:l="php artisan storage:link"

# telescope
alias ats="php artisan telescope"
alias ats:c="php artisan telescope:clear"
alias ats:i="php artisan telescope:install"
alias ats:pr="php artisan telescope:prune"
alias ats:p="php artisan telescope:publish"

# tinker
alias at="php artisan tinker"

# ui
alias aui="php artisan ui"
alias aui:a="php artisan ui:auth"

# up
alias au="php artisan up"

# vendor
alias avn="php artisan vendor"
alias avn:p="php artisan vendor:publish"

# view
alias av="php artisan view"
alias av:c="php artisan view:cache"
alias av:cl="php artisan view:clear"
